# triforce.service
[Unit]
Description=Triforce Multi-Agent AI Pipeline Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=@USER@
WorkingDirectory=@WORKDIR@
ExecStart=@NODE@ server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=triforce
Environment=NODE_ENV=production
EnvironmentFile=-@WORKDIR@/.env

[Install]
WantedBy=multi-user.target
